# sample.crontab
0 20 * * * bash ~/.ip-reporter/run.sh >> ~/.ip-reporter/run.log
